# roles/os-config/templates/disable-thp.service.j2
[Unit]
Description=disable-thp script

[Service]
ExecStart=/usr/local/bin/disable-thp

[Install]
WantedBy=multi-user.target
